# contrib/init/org.dogecoin.dogecoind.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.dogecoin.dogecoind</string>
	<key>ProgramArguments</key>
	<array>
<<<<<<< HEAD:contrib/init/org.dogecoin.dogecoind.plist
		<string>/usr/local/bin/dogecoind</string>
		<string>-daemon</string>
=======
		<string>/usr/local/bin/bitcoind</string>
>>>>>>> 1.21-dev:contrib/init/org.bitcoin.bitcoind.plist
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
